# integrationnet/supervisor.d/fixtures.conf
[program:fixtures]
user=stellar
command=/home/stellar/fixtures/load.sh
autostart=true
autorestart=false
priority=1000